# gradle/gradle/scripts/openbaton-nfvo.conf
# Open Baton init script

# Copyright (c) 2016 Open Baton (http://www.openbaton.org)
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
# 
#      http://www.apache.org/licenses/LICENSE-2.0
# 
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.


env OPENBATON_COMPONENT_NAME_FANCY="NFV Orchestrator (NFVO)"

description     "Open Baton ${OPENBATON_COMPONENT_NAME_FANCY}"
author          "Open Baton Team <info@openbaton.org>"

console output
start on (started networking)
stop on shutdown
kill timeout 60 
respawn limit 5 60	
limit core unlimited unlimited


env OPENBATON_COMPONENT_NAME=nfvo
env OPENBATON_COMPONENT_RELATIVE_SUBPATH=nfvo

# Expansion of variable like that does not work
#env OPENBATON_COMPONENT_CONFIG_FILE=/etc/openbaton/openbaton-${OPENBATON_COMPONENT_NAME}.properties
#env OPENBATON_COMPONENT_JAR_FILE_NAME=openbaton-${OPENBATON_COMPONENT_NAME}
#env OPENBATON_COMPONENT_JAR_INSTALL_DIR=/usr/lib/openbaton/${OPENBATON_COMPONENT_RELATIVE_SUBPATH}
#env OPENBATON_COMPONENT_TMP_FILE=/tmp/openbaton-${OPENBATON_COMPONENT_NAME}

env OPENBATON_COMPONENT_CONFIG_FILE=/etc/openbaton/openbaton-nfvo.properties
env OPENBATON_COMPONENT_JAR_FILE_NAME=openbaton-nfvo
env OPENBATON_COMPONENT_JAR_INSTALL_DIR=/usr/lib/openbaton/nfvo
env OPENBATON_COMPONENT_TMP_FILE=/tmp/openbaton-nfvo


###############
##### Main ####
###############

# Check the env is ready to start Open Baton
pre-start script
	OSTYPE=$(uname)

	####################
	#### Open Baton ####
	####################

	# Check if the Open Baton "${OPENBATON_COMPONENT_NAME}" is already running
    instance=$(ps aux | grep -v grep | grep "${OPENBATON_COMPONENT_JAR_FILE_NAME}" | grep "jar" | wc -l)
    if [ ${instance} -ne 0 ]; then
        logger -is -t "$UPSTART_JOB" "WARN: The Open Baton ${OPENBATON_COMPONENT_NAME_FANCY} is already running .."
        exit 1
    fi
	# Select which Open Baton "${OPENBATON_COMPONENT_NAME}" version to start
    openbaton_component_jar_file_name_versioned=$(ls -ct ${OPENBATON_COMPONENT_JAR_INSTALL_DIR} | grep "${OPENBATON_COMPONENT_JAR_FILE_NAME}" | grep "jar" | head -n 1)
	if [ "${openbaton_component_jar_file_name_versioned}" = "" ]; then
        logger -is -t "$UPSTART_JOB" "ERROR: The Open Baton ${OPENBATON_COMPONENT_NAME_FANCY} is not installed (jar-file not found)"
        exit 1
    fi
    echo "${openbaton_component_jar_file_name_versioned}" > "${OPENBATON_COMPONENT_TMP_FILE}"


	###################
	#### RabbitMQ #####
	###################

	# Check and eventually starts rabbitmq-server
    if [ "${OSTYPE}" = "Linux" ]; then
        instance=$(ps aux | grep -v grep | grep rabbitmq | grep server | wc -l)
        if [ ${instance} -eq 0 ]; then
            logger -is -t "$UPSTART_JOB" "WARN: rabbit is not running, let's try to start it .."
            #Start rabbitmq
			$(rabbitmq-server -detached)
            if (( $? != 0 )); then
                logger -is -t "$UPSTART_JOB" "ERROR: rabbitmq is not running properly (check the problem in /var/log/rabbitmq.log) "
                exit 1 
            fi
        fi
    elif [ "${OSTYPE}" = "Darwin" ]; then
        instance=$(ps aux | grep -v grep | grep rabbitmq | grep server | wc -l)
        if [ ${instance} -eq 0 ]; then
            logger -is -t "$UPSTART_JOB" "WARN: rabbitmq is not running, let's try to start it .."
            #Start rabbitmq
            $(rabbitmq-server -detached)
            if [ $? -ne 0 ]; then
                logger -is -t "$UPSTART_JOB" "ERROR: rabbitmq is not running properly (check the problem in /var/log/rabbitmq.log) "
                exit 1 
            fi
        fi
    fi
end script

# Start the Open Baton "${OPENBATON_COMPONENT_NAME}" as a service
script
	exec java -jar ${OPENBATON_COMPONENT_JAR_INSTALL_DIR}/$(cat ${OPENBATON_COMPONENT_TMP_FILE}) --spring.config.location=file:${OPENBATON_COMPONENT_CONFIG_FILE} --no-console
end script 

post-stop script
	rm -f ${OPENBATON_COMPONENT_TMP_FILE}
end script
